# demo2/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off";

[program:tomcat]
command=/usr/share/tomcat8/bin/catalina.sh run

